# systemd/nfs-auto-mount.timer
[Unit]
Description=Delay mount NFS after boot

[Timer]
OnBootSec=30sec
Unit=nfs-auto-mount.service

[Install]
WantedBy=timers.target
